# root/etc/services.d/youtube-dl-server/run
#!/usr/bin/with-contenv bash

exec \
        s6-setuidgid abc python -u /app/youtube-dl-server.py
